# cron/cron.d/backup
* * * * * thekashifmalik cd dotfiles && bin/backup
