# cron/cron.d/backup
* * * * * thekashifmalik touch ~/dotfiles/backup.`date -I`.log
