# cron/cron.d/backup
0 * * * * thekashifmalik cd dotfiles && bin/backup
