# cron/cron.d/restart-somememes-ingestion-backfill
*/15 * * * * thekashifmalik cd somememes && docker compose restart backend
